# supervisord.conf
; supervisor config file
; to be used with the default install only
; make sure to change paths if using global supervisor process 

[unix_http_server]
file=/tmp//supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/tmp/roombooking-supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/tmp/roombooking-supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/tmp/            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp//supervisor.sock ; use a unix:// URL  for a unix socket

[program:roombooking]
priority=0
command=%(here)s/virt/bin/python server.py
autostart=true
autorestart=true
directory=%(here)s/
redirect_stderr=true
redirect_stdout=true
